# templates/gunicorn-ini.in
#!/bin/sh

### BEGIN INIT INFO
# Provides:       drumcoder
# Required-Start: $local_fs $syslog
# Required-Stop:  $local_fs $syslog
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Gunicorn processes for drumcoder
### END INIT INFO

# see: http://drumcoder.co.uk/blog/2010/jul/30/simple-debian-init-script-gunicorn/


CONFDIR="/home/tjs/web/drumcoder"
NAME="drumcoder"
PORT=8001
WORKERS=1

USER=www-data
GUNICORN_RUN="gunicorn_django"
PID="/tmp/gunicorn_"$NAME".pid"
RETVAL=0

# source function library
. /lib/lsb/init-functions

start()
{
    echo "Starting $NAME."
    export PYTHONPATH=$CONFDIR;
    cd $CONFDIR/$NAME;
    su -c "$GUNICORN_RUN -b localhost:$PORT -n $NAME -w $WORKERS -p $PID -D settingslive.py" $USER && echo "OK" || echo "failed";
}

stop()
{
    echo "Stopping $NAME"
    kill -QUIT `cat $PID` && echo "OK" || echo "failed";
}

reload()
{
    echo "Reloading $NAME:"
    if [ -f $PID ]
    then kill -HUP `cat $PID` && echo "OK" || echo "failed";
    fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        reload
        ;;
    reload)
        reload
        ;;
    force-reload)
        stop && start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart}"
        RETVAL=1
esac
exit $RETVAL
